# config/crontab
15 16 * * * /rails/bin/rails runner "Rails.logger.warn 'This is a test cron job'"
