# config/crontab
15 16 * * * /rails/bin/rails runner "puts hello from uranus"
