# config/crontab
*/15 * * * * /rails/bin/cron-executor.sh bin/rails runner -e production 'puts User.first'
