# config/crontab
*/15 * * * * /rails/bin/cron-executor.sh ls -l
