# config/crontab
45 16 * * * /rails/bin/rails runner "puts 'hello from uranus'"
