# config/crontab
30 17 * * * cd /rails && bin/rails runner "puts 'hello from uranus'" > /tmp/uranus.log 2>&1
