# config/crontab
*/15 * * * * /rails/bin/cron_executor.sh ls -l
