# config/crontab
SHELL=/bin/bash

30 17 * * * cd /rails && bin/rails runner "puts 'hello from uranus'" > /proc/1/fd/1 2>&1
